# angel-cloud/angel-cloud-github-poller.timer
[Unit]
Description=Run Angel Cloud GitHub Poller every 15 minutes

[Timer]
OnBootSec=5min
OnUnitActiveSec=15min
Persistent=true

[Install]
WantedBy=timers.target
